# .config/autostart/sxhkd.desktop
[Desktop Entry]
Version=1.0
Name=sxhkd
Exec=sxhkd
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
