# .config/autostart/tilda.desktop
[Desktop Entry]
Version=1.0
Name=Tilda
Exec=sh -c "sleep 2; exec tilda --hidden"
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
