# .config/autostart/tilda.desktop
[Desktop Entry]
Version=1.0
Name=Tilda
Exec=tilda
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
